# systemd/remote-usb-server-attach.service.example
[Unit]
Description=Remote USB server — use devices shared by a client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust client IP, filters, and binary path.
ExecStartPre=/usr/local/bin/remote-usb server prepare
ExecStart=/usr/local/bin/remote-usb server 0.0.0.0 --client 192.168.1.10 --auto --match 14cd:1212
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
