# systemd/remote-usb-client.service
[Unit]
Description=Remote USB client — share local USB devices with a server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/usr/local/bin/remote-usb prepare
# Prefer --match VID:PID so keyboards/mice are not shared.
ExecStart=/usr/local/bin/remote-usb share 0.0.0.0 --port 3240 --auto
Restart=on-failure
RestartSec=2

# No authentication: restrict with firewall on this client.

[Install]
WantedBy=multi-user.target
