# systemd/remote-usb-server-attach.service.example
[Unit]
Description=Use remote USB devices from export host (auto-attach)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust host IP, filters, and binary path.
ExecStartPre=/usr/local/bin/remote-usb host 192.168.1.10 prepare
ExecStart=/usr/local/bin/remote-usb host 192.168.1.10 --auto --match 14cd:1212
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
