# systemd/remote-usb-client.service
[Unit]
Description=Remote USB export daemon (share local USB devices)
Documentation=https://wiki.archlinux.org/title/USB/IP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Install binary to a stable path, e.g. /usr/local/bin/remote-usb
ExecStartPre=/usr/local/bin/remote-usb prepare
# Prefer --match VID:PID so keyboards/mice are not exported.
ExecStart=/usr/local/bin/remote-usb server 0.0.0.0 --port 3240 --auto
Restart=on-failure
RestartSec=2

# No authentication: restrict with firewall.

[Install]
WantedBy=multi-user.target
