# systemd/remote-usb-server-attach.service.example
[Unit]
Description=Remote USB server — wait for clients to attach devices
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/usr/local/bin/remote-usb serve prepare
ExecStart=/usr/local/bin/remote-usb serve 0.0.0.0
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
